# salt/roots/apps/crontab
@reboot {{ project_dir }}/start_server.sh
