# install/salt/roots/apps/crontab
@reboot /home/users/ode_frontend/django_ode/start_server.sh
